# uwsgi.conf
# Emperor uWSGI upstart script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --emperor /etc/uwsgi/vassals

# vim: ft=upstart
